# root/etc/services.d/nepobot/run
#!/usr/bin/with-contenv bash
echo "Starting Bot...."
for (( ; ; )); do
	let i++
	npm start --prefix /app
done

exit $?
